# deploy/web3-university-api.service
[Unit]
Description=Web3 University API
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/web3-university-api
ExecStart=/usr/bin/node dist/main.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production
EnvironmentFile=/opt/web3-university-api/.env

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/web3-university-api

# 资源限制
LimitNOFILE=65536
MemoryMax=1G

[Install]
WantedBy=multi-user.target
